# etc/cron.d/update-void-list
@daily                          root    /bin/sh /usr/local/bin/void-zones-update.sh && service local_unbound restart

